# role/poll/files/poll.service
[Unit]
Description=Poll service

[Service]
User=poll
Group=poll
ExecStart=/usr/local/bin/wsgi.py

[Install]
WantedBy=multi-user.target