# roles/result/files/result.service
[Unit]
Description=Result service

[Service]
User=result
Group=result
ExecStart=node /usr/local/bin/server.js

[Install]
WantedBy=multi-user.target